# crontab
*/1  * * * * php  /var/www/cli/app/bootstrap.php main check 1
*/5  * * * * php  /var/www/cli/app/bootstrap.php main check 5
*/10  * * * * php  /var/www/cli/app/bootstrap.php main check 10